# tele1.service
[Unit]
Description=TELE1 Automated Telemetry Data Collection Service
Documentation=file:///home/tele/tele/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=tele
WorkingDirectory=/home/tele/tele

# Main service execution
ExecStart=/home/tele/tele/tele1.sh

# Timeout: 2 hours + small buffer
TimeoutStartSec=7500

# Stop behaviour
TimeoutStopSec=10
KillMode=process

# On completion or failure, power off the system
OnFailure=poweroff.target
OnSuccess=poweroff.target

# Do not restart, this is a one-shot run
Restart=no

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tele1

# Resource limits
MemoryLimit=512M
CPUQuota=80%

[Install]
WantedBy=multi-user.target
